# main/glibc/nscd.initd
#!/sbin/openrc-run

command="/usr/sbin/nscd"
description="Name Service Cache daemon"
pidfile="/var/run/nscd/nscd.pid"
